# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn sunshine:app -b localhost:5000
directory=/usr/src/app
